# package/debian/DEBIAN/postinst
#/bin/sh

systemctl enable mlx4_br --now