# package/debian/DEBIAN/prerm
#!/bin/sh

systemctl stop mlx4_br
systemctl disable mlx4_br